# debian/halpid.postrm
#!/bin/sh
set -e

case "$1" in
    purge)
        rm -rf /var/lib/halpid
        rm -rf /var/log/halpid
        rm -f /etc/halpid/halpid.conf
        ;;
    remove)
        systemctl stop halpid.service || true
        systemctl disable halpid.service || true
        ;;
esac

#DEBHELPER#
